# root/etc/services.d/caddy/run
#!/usr/bin/with-contenv sh

exec caddy -conf /etc/Caddyfile $CADDY_OPTIONS