# root/etc/services.d/docker-gen/finish
#!/usr/bin/with-contenv sh
s6-svscanctl -t /s6