# root/etc/services.d/docker-gen/run
#!/usr/bin/with-contenv sh

exec docker-gen -notify "s6-svc -1 /var/run/s6/services/caddy" -watch /etc/caddy-proxy/caddy.tmpl /etc/Caddyfile